# com.user.downloads-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.downloads-organizer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__INSTALL_DIR__/organize_downloads.py</string>
    <string>--downloads-dir</string>
    <string>__HOME__/Downloads</string>
    <string>--verbose</string>
    <string>--notify</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>WatchPaths</key>
  <array>
    <string>__HOME__/Downloads</string>
  </array>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/downloads-organizer/downloads_organizer.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/downloads-organizer/downloads_organizer.err.log</string>
</dict>
</plist>
